# simple.conf
[program:myserver]
#command=uwsgi --socket 172.31.33.200:8080 -w run:app
command=gunicorn run:app -w 4 -b 0.0.0.0:80

[supervisord]
logfile=/home/ubuntu/supervisord.log
loglevel=debug
user=root
